# supervisord.conf
[program:node]
directory=/srv/www
command=npm run demo
autostart=true
autorestart=true